# network/dyndns.timer
[Unit]
Description=Update public DNS entry for this network

[Timer]
OnBootSec=5min
OnCalendar=*:0/5

[Install]
WantedBy=basic.target
