# network/startup.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/startup.sh

[Install]
WantedBy=multi-user.target
